# recipes-packages/base/resize-rootfs/postinst
#!/bin/sh
chmod +x /etc/resize-rootfs/resize-rootfs.sh
systemctl enable resize-rootfs.service
